# files/service/crontab
0 4 * * * root /usr/odk/purge-forms.sh
0 3 * * * root /usr/odk/run-analytics.sh
0 1 * * 0 root /usr/odk/reap-sessions.sh
